# analyst_gpt/crontab
*/10 * * * * cd /app && python -m tasks.run_summariser
